# roles/macos/templates/com.pedro.dotfiles.maintenance.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pedro.dotfiles.maintenance</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>{{ ansible_facts.env.HOME }}/Developer/dotfiles/scripts/maintenance.sh</string>
    </array>
    
    <!-- Run weekly on Sundays at 10:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Also run on load if missed (e.g., Mac was asleep) -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>{{ ansible_facts.env.HOME }}/.local/log/maintenance-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{ ansible_facts.env.HOME }}/.local/log/maintenance-stderr.log</string>
</dict>
</plist>
